# vsm.service
# deploy/vsm.service
# ────────────────────
# systemd service for VSM — Vehicle Safety Monitor
#
# Install on Raspberry Pi:
#   1. Copy this file:
#        sudo cp deploy/vsm.service /etc/systemd/system/vsm.service
#
#   2. Edit the paths below to match your installation:
#        User=        → your Pi username (default: pi)
#        WorkingDirectory= → absolute path to vsm/
#        ExecStart=   → path to your venv python + main.py
#
#   3. Enable and start:
#        sudo systemctl daemon-reload
#        sudo systemctl enable vsm
#        sudo systemctl start vsm
#
#   4. Check status / logs:
#        sudo systemctl status vsm
#        journalctl -u vsm -f
#
#   5. Stop:
#        sudo systemctl stop vsm

[Unit]
Description=VSM — Vehicle Safety Monitor
Documentation=https://github.com/yourorg/vsm
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/vsm

# Activate virtualenv and run
ExecStart=/home/pi/vsm/.venv/bin/python main.py

# Restart automatically if it crashes
Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=60s
StartLimitBurst=5

# Environment — override config path if needed
Environment="PYTHONUNBUFFERED=1"
Environment="VSM_CONFIG=/home/pi/vsm/config.yaml"

# Give time for camera / GPS to initialise on boot
ExecStartPre=/bin/sleep 5

# Graceful stop — sends SIGTERM, waits 10s, then SIGKILL
TimeoutStopSec=10s
KillMode=mixed

# Logging — output goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vsm

[Install]
WantedBy=multi-user.target
